# the-butler.service
[Unit]
Description=The Butler Discord Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=butlerbot
Group=butlerbot
WorkingDirectory=/opt/the-butler/app
EnvironmentFile=/opt/the-butler/app/.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/the-butler/app/.venv/bin/python /opt/the-butler/app/main.py
Restart=always
RestartSec=10
UMask=0077

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/the-butler/data /opt/the-butler/logs
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
SystemCallArchitectures=native

StandardOutput=append:/opt/the-butler/logs/the-butler.log
StandardError=append:/opt/the-butler/logs/the-butler.err.log

[Install]
WantedBy=multi-user.target
